# etc/crontabs/root
25 * * * * tail -n 7 /tmp/log/agent.log > /tmp/log/agent_.log; cat /tmp/log/agent_.log > /tmp/log/agent.log; rm /tmp/log/agent_.log
